# he6-tunnel.container
[Container]
Name=he-tunnel
Image=localhost/he-tunnel
Network=host
Capabilities=NET_ADMIN
